# .config/zsh/profile.d/linuxbrew
export LINUXBREW_HOME="${HOME}/.linuxbrew"
export PATH="${LINUXBREW_HOME}/bin:${LINUXBREW_HOME}/sbin:$PATH"
export MANPATH="${LINUXBREW_HOME}/share/man:$MANPATH"
export INFOPATH="${LINUXBREW_HOME}/share/info:$INFOPATH"
export XDG_DATA_DIRS="${LINUXBREW_HOME}/share:${XDG_DATA_DIRS}"
